# supervisord.conf
; Cwt Academy — Production Queue Worker Configuration
; Install: sudo apt install supervisor
; Copy to: /etc/supervisor/conf.d/cwt-academy-worker.conf
; Reload: sudo supervisorctl reread && sudo supervisorctl update

[program:cwt-academy-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/cwt-academy/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --queue=default
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cwt-academy-worker.log
stopwaitsecs=3600
environment=HOME="/var/www"

[program:cwt-academy-scheduler]
process_name=%(program_name)s
command=php /var/www/cwt-academy/artisan schedule:run
autostart=false
autorestart=false
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/supervisor/cwt-academy-scheduler.log
